# i2c-tools/build-deb/DEBIAN/prerm
#!/bin/sh

set -ex

svc=fan-ctrl

systemctl stop ${svc}.service
systemctl disable ${svc}.service

exit 0
